# server/gunicorn.socket
[Unit]
Description=Galaxy content site socket
[Socket]
ListenStream=/run/ga-site.sock
[Install]
WantedBy=sockets.target
