# deploy/ansible/files/webapp.socket
[Unit]
Description=Galaxy content site socket
[Socket]
ListenStream=/run/webapp.sock
[Install]
WantedBy=sockets.target
